# debpack/DEBIAN/postinst
#!/bin/bash
# Activate 9590
# A simple listener on Port 9590 for documentation examples
/usr/bin/systemctl daemon-reload
/usr/bin/systemctl enable 9590
/usr/bin/systemctl start 9590
