# debpack/DEBIAN/postinst
#!/bin/bash
# Activate 9590
# A simple listener on Port 9590 for documentation examples
/usr/sbin/update-rc.d 9590 defaults
/usr/sbin/update-rc.d 9590 enable
/etc/init.d/9590 start
